# files/system/etc/systemd/system/adw-gtk3-update.service
[Unit]
Description=Update adw-gtk3 theme

[Service]
Type=oneshot
ExecStart=/usr/bin/adw-gtk3-update

[Install]
WantedBy=multi-user.target